# examples/postdare-go.service
[Unit]
Description=Postdare Go release console
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/postdare-go
Environment=POSTDARE_GO_DATA_DIR=/data/postdare-go
ExecStart=/opt/postdare-go/postdare-go
Restart=always
RestartSec=3
User=postdare-go
Group=postdare-go

[Install]
WantedBy=multi-user.target
